# scripts/crontab
0 2 * * * cd /path/to/your/project/root && ./scripts/deploy_by_pull.sh
